# src/zfs-esp-sync@.service
[Unit]
Description=Sync ZFS ESP

[Service]
ExecStart=/usr/bin/zfs-esp-sync "%i"
