# src/zfs-esp-sync@.service
[Unit]
Description=Sync ZFS ESP

[Install]
WantedBy=basic.target

[Service]
ExecStart=/usr/bin/zfs-esp-sync-watch "%i"
