# src/zfs-esp-sync@.path
[Unit]
Description=Sync ZFS ESP

[Install]
WantedBy=basic.target

[Path]
PathChanged=/boot/initramfs-%i.img
PathChanged=/boot/initramfs-%i-fallback.img
PathChanged=/boot/intel-ucode.img
