# src/zfs-esp-sync@.path
[Unit]
Description=Sync ZFS ESP

[Path]
PathChanged=/boot/initramfs-%i.img
PathChanged=/boot/initramfs-%i-fallback.img
